# scripts/mac/com.local.scheduled-shutdown-bdt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.local.scheduled-shutdown-bdt</string>
	<key>Comment</key>
	<string>Shut down Mac daily at 2:00 AM Asia/Dhaka (independent of HANOON)</string>
	<key>ProgramArguments</key>
	<array>
		<string>/sbin/shutdown</string>
		<string>-h</string>
		<string>now</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>2</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/var/log/scheduled-shutdown-bdt.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/scheduled-shutdown-bdt.log</string>
</dict>
</plist>
